# systemd/scionlab.target
[Unit]
Description=SCIONLab Service

[Install]
WantedBy=multi-user.target
