# debian/scripts/scionlab.prerm
#!/usr/bin/env bash
set -e

if [ -f /.dockerenv ]; then
  exit 0
fi

case "$1" in
    # Stop services before pulling the rug from under systemd.
    upgrade)
        systemctl stop scionlab.target
        ;;
    remove)
        systemctl disable scionlab.target
        systemctl stop scionlab.target
        ;;
    *)
        ;;
esac
