# debian/scripts/scionlab.prerm
#!/usr/bin/env bash
set -e

if [ "$1" = remove ] && [ ! -f /.dockerenv ]; then
    systemctl disable scionlab.target
    systemctl stop scionlab.target
fi
